# crontab
# Run scraper every day at 3:00 AM UTC
0 3 * * * bun ./scrape/index.js >> /dev/stdout 2>> /dev/stderr