# crontab
# Scrape for data every 45 minutes
*/45 * * * * cd /tourRanking && /usr/local/bin/bun scrape
