# crontab
# Scrape for data every 45 minutes
*/45 * * * * cd /tourRanking && bun scrape >> /dev/stdout 2>> /dev/stderr
